# mongodb/roles/mongod/templates/mongod.service.j2
[Unit]
Description=High-performance, schema-free document-oriented database
After=network.target

[Service]
User=root
Type=forking
ExecStart=/usr/bin/mongod --quiet --maxConns=60000 --config /etc/mongod-{{ inventory_hostname }}.conf

[Install]
WantedBy=multi-user.target